# source/usr/local/emhttp/plugins/dwdvm/dwdvm.cron
# Data Volume Monitor Alarms Script:
* * * * * /etc/vnstat/dvm-alarms.sh > /dev/null 2>&1
