# source/usr/local/emhttp/plugins/dwdvm/dwdvm.cron
# Data Volume Monitor Custom Alarms Script:
* * * * * /etc/vnstat/dvm-custom.sh > /dev/null 2>&1
